# nextcloud-cron/crontabs/www-data
*/5 * * * * php -f /var/www/html/cron.php
*/5 * * * * php /var/www/html/occ preview:pre-generate
